# tests/service/sshd/run
#!/bin/sh

ssh-keygen -A
exec 2>&1
exec /usr/sbin/sshd -D -e
